# config/ruby_everyday_init.conf
description "Ruby Everyday"

start on runlevel [23]
stop on shutdown

exec sudo sh -c "bash /home/xshirowx/apps/ruby-everyday/current/config/unicorn_init.sh start"

respawn